# systemd/user/ikanban-proxy.service
[Unit]
Description=IKanban Task Management System (Proxy)
After=network.target graphical-session.target

[Service]
Type=simple
Environment="PATH=%h/.bun/bin:%h/.local/bin:%h/.nvm/versions/node/v24.12.0/bin:/usr/local/bin:/usr/bin"
ExecStart=/usr/bin/zsh -lc 'source %h/.zshrc; set -a; source %h/.token; set +a;HTTPS_PROXY=http://localhost:1080 %h/.opencode/bin/opencode serve --port 5000 --cors https://isomoes.github.io'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
